# SOS-Firewall.service
[Unit]
Description=Start SOS-Firewall

[Service]
ExecStart=/usr/bin/Start-SOS-Firewall

[Install]
WantedBy=multi-user.target 